# SOURCES/etc/cron.d/thruk_logcache
*/5 * * * * root /srv/rgm/thruk/script/thruk -a logcache update &> /dev/null
